# rpm/src/main/resources/exomiser.cron
*/15 * * * *   tomcat   bash /var/lib/exomiser/scripts/exomiser_cron.sh
